# imageroot/systemd/user/certificate-exporter.path
[Unit]
Description=Monitor acme.json file for changes

[Path]
PathChanged=%S/state/acme/acme.json
Unit=certificate-exporter.service
